# dns-fallback.service
[Unit]
Description=DNS Fallback Pi-hole Proxy Service
After=network.target pihole-FTL.service unbound.service
Wants=pihole-FTL.service unbound.service

[Service]
Type=simple
WorkingDirectory=/opt/dns-fallback
ExecStartPre=/bin/rm -f /var/run/dns-fallback.pid
ExecStart=/opt/dns-fallback/venv/bin/python3 /opt/dns-fallback/dns_fallback_proxy.py
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=append:/var/log/dns-fallback.log
StandardError=append:/var/log/dns-fallback.log
User=root
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/log /opt/dns-fallback /var/run

[Install]
WantedBy=multi-user.target
