# dns-fallback-dashboard.service
[Unit]
Description=DNS Fallback Pi-hole Dashboard Service
After=network.target dns-fallback.service

[Service]
Type=simple
User=root
# If you want to run as a dedicated user (recommended for security):
# User=dnsfallback
# Group=dnsfallback
WorkingDirectory=/opt/dns-fallback
ExecStart=/usr/bin/python3 /opt/dns-fallback/dns_fallback_dashboard.py
Restart=on-failure
RestartSec=5s
StandardOutput=file:/var/log/dns-fallback_dashboard.log
StandardError=file:/var/log/dns-fallback_dashboard.log
# Alternative for logging: use journald
# StandardOutput=journal
# StandardError=journal

[Install]
WantedBy=multi-user.target
